# crons/db-backup.cron
# managed-by:.github
# Exemple de tâche cron : sauvegarde DB tous les jours à 3h du matin
0 3 * * * /home/ichisamafr/dev/.github.git/scripts/db/backup.sh >> /var/log/db-backup.log 2>&1
